# energy_policy.service
[Unit]
Description=Set default cpu power policy
[Service]
Type=oneshot
ExecStart=/usr/local/bin/cpu 1
[Install]
WantedBy=multi-user.target
